# checkra1n.service
[Unit]
Description = Launch checkrain on launch
 
[Service]
ExecStart=/opt/microrain/Xcheckra1n
 
[Install]
WantedBy=multi-user.target